# databases/redis/files/redis.in
#!/bin/sh
#
# $FreeBSD$
#
# PROVIDE: redis
# REQUIRE: LOGIN
# BEFORE:  securelevel
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf to run redis:
#
# redis_enable (bool):       Set it to "YES" to enable redis.
#                            Default is "NO".
# redis_profiles (strings):  Example - "foo bar"
#                            Default is empty.

. /etc/rc.subr

name="redis"
rcvar=redis_enable

load_rc_config ${name}

: ${redis_enable="NO"}
: ${redis_user="redis"}

command="%%PREFIX%%/bin/redis-server"


_profile_exists() {
	for _p in ${redis_profiles}; do
		[ "${_p}" = "$1" ] && return 1;
	done
	return 0
}

if [ $# -eq 2 ]; then
	_profile=$2
	_profile_exists $_profile
	_exists=$?
	[ ${_exists} -ne 1 ] && {
		echo "`basename /usr/local/etc/rc.d/redis`: no '$2' in 'redis_profiles'"
		exit 1
	};
	echo "-- Profile: ${_profile} --"
	command_args="/usr/local/etc/${name}-${_profile}.conf"
	pidfile="/var/run/redis/${_profile}.pid"
elif [ -n "${redis_profiles}" ]; then
	_swap=$*; shift; _profiles=$*
	_profiles=${_profiles:-${redis_profiles}}
	set -- ${_swap}
	for _profile in ${_profiles}; do
		/usr/local/etc/rc.d/redis $1 ${_profile}
	done
	exit 0
else
	command_args="/usr/local/etc/${name}.conf"
	pidfile="/var/run/redis/${name}.pid"
fi


run_rc_command "$1"
